# ops/rotate_logs.timer
# rotate_logs.timer
[Unit]
Description=Daily rotate logs

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target